# deploy/systemd/cc-switch-web.service.example
[Unit]
Description=CC Switch Web local Rust service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ccswitch
Group=ccswitch
WorkingDirectory=/opt/cc-switch-web
Environment=HOME=/home/ccswitch
Environment=CC_SWITCH_WEB_HOST=127.0.0.1
Environment=CC_SWITCH_WEB_PORT=8788
ExecStart=/opt/cc-switch-web/cc-switch-web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
